# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
minfds=1024
minprocs=200
autorestart=true

[program:php]
command=%(ENV_SUPERVISOR_PHP_COMMAND)s
directory=/var/www/html
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:inertia-ssr]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/html/artisan inertia:start-ssr
directory=/var/www/html
user=sail
autostart=true
autorestart=true
startretries=5
stdout_logfile=/var/www/html/storage/logs/inertia_ssr_stdout.log
stderr_logfile=/var/www/html/storage/logs/inertia_ssr_stderr.log
stopsignal=QUIT
stopwaitsecs=10
